# roles/neutron-data/templates/etc/cron.daily/cleanup-neutron-interfaces
#!/bin/sh

# Remove logs for Neutron TAP interfaces
find /var/log/upstart/ -name network-interface-tap\*.log -mtime +1 -delete
