# roles/iptables/templates/etc/network/if-up.d/iptables
#!/bin/bash
iptables-restore < /etc/network/iptables-firewall
