# roles/swift-common/templates/etc/update-motd.d/99-swift-motd
#!/bin/sh

echo
echo "WARNING Swift node -- do not replace disks without consulting Swift specific procedures"
echo
